# service/reload_plex.service
[Unit]
Description=Auto reload PLEX device

[Service]
ExecStart=/usr/local/bin/reload_plex.sh
Restart=always

[Install]
WantedBy=multi-user.target
